# home/.config/systemd/user/extension-prefs-monitor-laptop.service
[Unit]
Description=GNOME Extension Preferences Monitor
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/extension-prefs-monitor.sh %h/Projects/setup-scripts/fedora/laptop-extensions-prefs.conf
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
